# bootscripts/s6/services/hostname/up
backtick -En hostname { cat /etc/hostname } s6-hostname $hostname
